# package/system/systemd/files/run.mount
[Unit]
Description=Run Directory
DefaultDependencies=no
Before=sysinit.target
ConditionPathExists=!/run

[Mount]
What=tmpfs
Where=/run
Type=tmpfs
Options=mode=755,nosuid,nodev,size=10%

[Install]
WantedBy=local-fs.target
